# config/voice-assistant.service
[Unit]
Description=TrailCurrent Peregrine Voice Assistant
After=network.target sound.target genie-server.service
Wants=genie-server.service

[Service]
Type=simple
User=trailcurrent
Group=audio
WorkingDirectory=/home/trailcurrent
ExecStartPre=/bin/sleep 5
ExecStart=/home/trailcurrent/assistant-env/bin/python3 /home/trailcurrent/assistant.py
Restart=on-failure
RestartSec=10

# Environment
Environment=HOME=/home/trailcurrent
Environment=WAKE_MODEL_PATH=/home/trailcurrent/models/hey_peregrine.onnx
Environment=PATH=/home/trailcurrent/assistant-env/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1
# Site-specific config (MQTT, thresholds, etc.) — never overwritten by deploys.
EnvironmentFile=-/home/trailcurrent/assistant.env

# Hardening
ProtectSystem=strict
ReadWritePaths=/home/trailcurrent /tmp
ProtectHome=tmpfs
BindPaths=/home/trailcurrent
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
