# image_build/files/systemd/power-save-hw.service
[Unit]
Description=Disable unused QCS6490 hardware for power savings
After=multi-user.target plymouth-quit.service

[Service]
Type=oneshot
ExecStart=/bin/bash -c '\
  for gpu in /sys/class/devfreq/*gpu* /sys/class/devfreq/*3d00000*; do \
    [ -f "$gpu/governor" ] && echo powersave > "$gpu/governor"; \
    min=$(awk "{print \$1}" "$gpu/available_frequencies" 2>/dev/null); \
    [ -n "$min" ] && echo "$min" > "$gpu/max_freq" && echo "$min" > "$gpu/min_freq"; \
  done; \
  for drv in msm_dsi msm_dp msm_mdss camss; do \
    d="/sys/bus/platform/drivers/$drv"; [ -d "$d" ] && \
    for dev in "$d"/*/; do n=$(basename "$dev"); \
      [ "$n" != module ] && [ "$n" != uevent ] && echo "$n" > "$d/unbind" 2>/dev/null; \
    done; \
  done; \
  rfkill block wifi 2>/dev/null || true; \
  for dev in /sys/bus/usb/devices/*/power/control; do \
    dp=$(dirname "$dev"); skip=false; \
    for ifc in "$dp"/*:*/bInterfaceClass; do \
      [ -f "$ifc" ] || continue; \
      class=$(cat "$ifc" 2>/dev/null | tr -d "[:space:]"); \
      case "$class" in 01|03) skip=true; break;; esac; \
    done; \
    $skip && echo on > "$dev" 2>/dev/null || echo auto > "$dev" 2>/dev/null; \
  done; \
  true'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
